# deploy/iris-agent.service
[Unit]
Description=Iris LiveKit agent worker
# Wait for the LiveKit server containers to be up before starting; if the
# agent starts first, it'll just retry until LiveKit responds, but waiting
# avoids restart-loop noise on droplet reboot.
After=network.target docker.service
Wants=docker.service
# Flip the unit to "failed" state after 3 quick restart attempts in 60s.
# Combined with the prewarm-time _smoke_test_constructors() in iris_agent.py,
# a bad deploy (e.g., the keywords/keyterm bug that hit 2026-06-16 to 06-20)
# causes the worker to fail prewarm 3x and then systemctl status shows red.
# deploy.bat's status-check step and any external uptime monitor can see it.
# Without StartLimit, Restart=always would just retry the broken process
# forever, hiding the failure from operational checks.
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
User=iris
Group=iris
WorkingDirectory=/opt/iris-backend/agent
EnvironmentFile=/opt/iris-backend/agent/.env
# Run via the venv's interpreter so we don't need to source activate.
# `start` (not `dev`) skips the file-watcher reload that `dev` enables.
ExecStart=/opt/iris-backend/agent/.venv/bin/python /opt/iris-backend/agent/iris_agent.py start
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Hardening (matches iris-backend.service)
NoNewPrivileges=true
ProtectSystem=strict
# Writable paths:
#  - /home/iris/.cache — Hugging Face model downloads
#  - /opt/iris-backend/agent — venv + transient artifacts
#  - /opt/iris-backend/recordings — per-call transcript JSONs the agent
#    writes from its shutdown callback (Egress also writes OGGs here from
#    inside its container — that path bypasses systemd hardening)
ReadWritePaths=/home/iris/.cache /opt/iris-backend/agent /opt/iris-backend/recordings
PrivateTmp=true

# /run/iris/active_calls — call-active flag dir read by the diarize
# watcher daemon (iris-diarize-watcher.service). Group-writable so the
# watcher (also User=iris, Group=iris) can prune stale flags. Preserve=yes
# keeps the directory across restarts of THIS unit so the watcher's
# in-flight diarize doesn't lose its coordination signal if the agent
# restarts mid-call. Mode 0775 = owner/group rwx, others rx.
RuntimeDirectory=iris/active_calls
RuntimeDirectoryMode=0775
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
