# deploy/iris-diarize-watcher.service
[Unit]
Description=Iris per-call diarize watcher
# The watcher tolerates the agent being absent (it'll just see no call
# flags and run diarize freely). Order helps log readability.
After=network.target iris-agent.service
Wants=iris-agent.service

[Service]
Type=simple
User=iris
Group=iris
WorkingDirectory=/opt/iris-backend/tools/diarize

# Use the diarize venv — separate from the agent venv (Whisper +
# pyannote vs LiveKit Agents). The watcher itself only uses stdlib,
# but it execs the venv's python for the diarize batch child.
ExecStart=/opt/iris-backend/tools/diarize/.venv/bin/python /opt/iris-backend/tools/diarize/diarize_watcher.py

# Restart on crash with a small backoff so a persistent error (e.g.,
# corrupt OGG that triggers a Whisper crash) doesn't hot-loop.
Restart=on-failure
RestartSec=10s

# Shared runtime dir with the agent. See iris-agent.service for
# RuntimeDirectory= rationale. Both services declare the same dir;
# systemd creates it once and both can read/write (matching User/Group).
RuntimeDirectory=iris/active_calls
RuntimeDirectoryMode=0775
RuntimeDirectoryPreserve=yes

# NO cgroup memory/CPU caps. The previous version had MemoryMax=256M
# and CPUQuota=10% intended for the tiny watcher daemon — but systemd
# cgroup limits apply to the WHOLE service cgroup including child
# processes. The diarize batch needs ~3-4 GB for Whisper large-v3 int8
# and gets ~50%+ of one CPU even at nice -n 19; both were strangled
# by those limits. Verified 2026-06-16: the diarize child got stuck in
# mem_cgroup_handle_over_high (kernel throttle) for 2:45 with no
# forward progress before this fix.
#
# Resource management is handled OUTSIDE the cgroup:
#   - nice -n 19 / ionice -c 3 on the diarize child (in diarize_watcher.py)
#     for CPU/IO yielding when other work is ready.
#   - Explicit SIGKILL by the watcher when /proc/meminfo's MemAvailable
#     drops below MEM_PAUSE_THRESHOLD_MB during a live call.
# Both mechanisms are unaffected by cgroup-level limits.

# Hardening (matches iris-agent.service)
NoNewPrivileges=true
ProtectSystem=strict
# Writable paths:
#  - /opt/iris-backend/recordings — diarize batch writes transcribed
#    JSONs here under transcribed/
#  - /opt/iris-backend/tools/diarize — venv + any cached models
#  - /home/iris/.cache — Hugging Face pyannote model cache
ReadWritePaths=/opt/iris-backend/recordings /opt/iris-backend/tools/diarize /home/iris/.cache
PrivateTmp=true

# Environment overrides (uncomment to tune; defaults are calibrated
# for the 8 GB droplet):
#Environment=MEM_PAUSE_THRESHOLD_MB=1500
#Environment=CALL_POLL_INTERVAL_S=2.0
#Environment=OGG_POLL_INTERVAL_S=10.0
#Environment=IRIS_RECORDINGS_DIR=/opt/iris-backend/recordings
#Environment=IRIS_CALL_ACTIVE_DIR=/run/iris/active_calls
#Environment=DIARIZE_VENV=/opt/iris-backend/tools/diarize/.venv

StandardOutput=journal
StandardError=journal
SyslogIdentifier=iris-diarize-watcher

[Install]
WantedBy=multi-user.target
