# deploy/iris-agent.service
[Unit]
Description=Iris LiveKit agent worker
# Wait for the LiveKit server containers to be up before starting; if the
# agent starts first, it'll just retry until LiveKit responds, but waiting
# avoids restart-loop noise on droplet reboot.
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=iris
Group=iris
WorkingDirectory=/opt/iris-backend/agent
EnvironmentFile=/opt/iris-backend/agent/.env
# Run via the venv's interpreter so we don't need to source activate.
# `start` (not `dev`) skips the file-watcher reload that `dev` enables.
ExecStart=/opt/iris-backend/agent/.venv/bin/python /opt/iris-backend/agent/iris_agent.py start
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Hardening (matches iris-backend.service)
NoNewPrivileges=true
ProtectSystem=strict
# Writable paths:
#  - /home/iris/.cache — Hugging Face model downloads
#  - /opt/iris-backend/agent — venv + transient artifacts
#  - /opt/iris-backend/recordings — per-call transcript JSONs the agent
#    writes from its shutdown callback (Egress also writes OGGs here from
#    inside its container — that path bypasses systemd hardening)
ReadWritePaths=/home/iris/.cache /opt/iris-backend/agent /opt/iris-backend/recordings
PrivateTmp=true

[Install]
WantedBy=multi-user.target
