# deploy/iris-backend.service
[Unit]
Description=Iris (Lighthouse Inn AI reservation agent backend)
After=network.target

[Service]
Type=simple
User=iris
Group=iris
WorkingDirectory=/opt/iris-backend/backend
EnvironmentFile=/opt/iris-backend/backend/.env
ExecStart=/opt/iris-backend/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
# /opt/iris-backend/backend/data — SQLite DB lives here.
# /opt/iris-backend/recordings — the dashboard at /iris/ writes merged
#   stereo OGGs (merged_{call_id}.ogg) and Claude summary sidecars
#   (summary_{call_id}.json) into this directory alongside the agent's
#   transcripts and per-track OGGs.
ReadWritePaths=/opt/iris-backend/backend/data /opt/iris-backend/recordings
PrivateTmp=true

[Install]
WantedBy=multi-user.target
